# service/binarystore-migration.service
[Unit]
Description=Migrate Binary Stroage Data

[Service]
Type=oneshot
ExecStart=blobtool -m

