# service/binarystore-migration.service
[Unit]
Description=Migrate Binary Storage Data

[Service]
Type=oneshot
ExecStart=/usr/bin/blobtool -m

